# roles/hdfs/templates/etc/profile.d/hadoop.sh.j2
# {{ ansible_managed }}

export PATH="$PATH:{{ hadoop_home }}/bin:{{ hadoop_home }}/sbin"
